# scripts/power_diagnostics/sf-power-marker.service
[Unit]
Description=Smart Factory power-loss marker and boot diagnostics
DefaultDependencies=no
After=local-fs.target systemd-journald.service
Before=multi-user.target shutdown.target reboot.target halt.target
Conflicts=shutdown.target reboot.target halt.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/sbin/sf-power-diagnostics boot
ExecStop=/usr/local/sbin/sf-power-diagnostics shutdown
TimeoutStartSec=30
TimeoutStopSec=10

[Install]
WantedBy=multi-user.target
