# scripts/power_diagnostics/sf-power-snapshot.service
[Unit]
Description=Smart Factory periodic power diagnostics snapshot

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/sf-power-diagnostics snapshot
